# buildroot/board/rockchip/rv1126_rv1109/fs-overlay-sysv/etc/profile.d/RkEnv.sh
export HOME=
export HOME_OEM=/oem
export PATH=$PATH:$HOME:$HOME/bin:$HOME/usr/bin:$HOME/sbin:$HOME/usr/sbin:$HOME_OEM/bin:$HOME_OEM/usr/bin:$HOME_OEM/sbin:$HOME_OEM/usr/sbin
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:$HOME/usrlibs:$HOME/usr/lib:$HOME/lib:$HOME_OEM/usr/lib:$HOME_OEM/lib

